# rootfs/etc/services.d/opendkim/run
#!/usr/bin/with-contenv sh

mkdir -p /run/opendkim
exec /usr/sbin/opendkim -f -l -p inet:8891
